# sound/SerenityAudio/serenityaudio.service
[Unit]
Description=Flaming Lotus Girls 2019 Serenity SoundScape Server

[Service]
Type=simple
User=flaming
WorkingDirectory=/home/flaming/Serenity/sound/SerenityAudio
ExecStart=/home/flaming/Serenity/sound/SerenityAudio/saplay
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
